# sentry/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
[program:{{ service_name }}]
directory=/srv/sentry/run
environment=PATH="/srv/sentry/bin";SENTRY_CONF="/etc/sentry/sentry.conf.py"
user=sentry
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/sentry/logs/{{ service.name }}.log
{%- if service.name == 'web' %}
command=/srv/sentry/bin/sentry --config=/etc/sentry/sentry.conf.py run web
{%- elif service.name == 'worker' %}
command=/srv/sentry/bin/sentry --config=/etc/sentry/sentry.conf.py run worker -l WARNING
{%- elif service.name == 'cron' %}
command=/srv/sentry/bin/sentry --config=/etc/sentry/sentry.conf.py run cron
{%- endif %}
